# recipes-graphics/x11-common/xserver-nodm-init/xserver-nodm.path
[Unit]
Description="DRI graphics interface gating unit"

[Path]
PathExists=/dev/dri/card0

[Install]
WantedBy=graphical.target
